# docker/etc/services.d/djmngr/run
#!/usr/bin/execlineb -P
/app/bin/server-docker
